# supervisord.conf
[supervisord]
nodaemon=true

[program:suricata]
command=/bin/bash -c "exec suricata -i $(route | grep default | awk '{ print $8 }')"
autorestart=true

[program:p0f]
directory=/opt/p0f/
command=/bin/bash -c "exec /opt/p0f/p0f -i $(route | grep default | awk '{ print $8 }') -j -o /data/suricata/log/p0f.json"
autorestart=true
