# supervisord.conf
[supervisord]
nodaemon=true

[program:suricata]
command=/bin/bash -c "exec suricata -i $(route | grep default | awk '{ print $8 }')"
autorestart=true
